# crontab
0 23 * * * root /opt/app/src/run.sh >> /var/log/cron.log 2>&1
